# provision/service/beanstalkd.conf
[program:beanstalkd]
command=/usr/bin/beanstalkd -l 127.0.0.1 -p 11300
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true