# provision/service/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
