# provision/service/nginx.conf
[program:nginx]
command=/usr/sbin/nginx
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true