# provision/service/php-fpm.conf
[program:php5-fpm]
command=/usr/sbin/php-fpm7.0 -c /etc/php/7.0/fpm
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
