# provision/service/queue.conf
[program:queue]
command=php /share/artisan queue:listen --tries=1
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
